# collect_health.service
[Unit]
Description=Collect Raspberry Pi Health Data

[Service]
Type=oneshot
ExecStart=/usr/local/bin/collect_health.py
