# netcheck.service
[Unit]
Description=Network Connectivity Check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/netcheck.sh
