# secscan.service
[Unit]
Description=Security Log Collector

[Service]
Type=oneshot
ExecStart=/usr/local/bin/secscan.sh
