# secscan.timer
[Unit]
Description=Run Security Scan every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Unit=secscan.service

[Install]
WantedBy=timers.target
